# configs/supervisord/example.com.conf
[program:xmpp-bot]
command = /home/user/opt/microblog/env/bin/bot /home/user/opt/microblog/configs/bot/example.cfg
